# roles/base/templates/apt-refresh.service.j2
[Unit]
Description=Refresh workstation apt package lists
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/apt-refresh
RuntimeMaxSec=10m
