# roles/base/templates/managed-package-updates.timer.j2
[Unit]
Description=Schedule daily managed baseline package upgrades

[Timer]
OnCalendar={{ base_managed_package_updates_timer_on_calendar }}
Persistent=true

[Install]
WantedBy=timers.target
