# roles/base/templates/managed-package-updates.service.j2
[Unit]
Description=Upgrade installed managed baseline packages
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/upgrade-installed-apt-packages --label managed-baseline --list-file /etc/ansible/managed-package-updates.list
